# crontab
# START CRON JOB
0 * * * * /usr/local/bin/python3 /app/DNSUpdater.py
# END CRON JOB